# supervisor/conf.d/prana.conf
[program:prana]
command=/usr/local/bin/prana.sh
autostart=true